# supervisor/conf.d/relayctl.conf
[program:relayctl]
#This _may_ fail once or twice waiting for pigpiod to come up
command=/usr/bin/python3 /usr/local/sbin/relayctl.py /etc/relayctl/relayctl.conf.json
autorestart=true
user=relayctl